# configs/borg/borgbackup-home.service
[Unit]
Description=Home Backup using BorgBackup

[Service]
Type=oneshot
ExecStart=/opt/borg/home-borgbackup.sh
